# .local/share/systemd/user/proton-update.timer
[Unit]
Description=Proton Updater

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
