# .local/share/systemd/user/spicetify-update.timer
[Unit]
Description=Spicetify Updater

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
