# .local/share/systemd/user/who-dl.timer
[Unit]
Description=Doctor Who Episode Scraper

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
